# deployment/supervisor-web.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /code/uwsgi.ini --uid web

[program:nginx]
command = /usr/sbin/nginx

[program:cron]
command = /usr/sbin/cron -f
